# usr/lib/systemd/system/make-swapfile-to-hibernate.service
[Unit]
Description=Create swapfile for hibernate
Before=systemd-hibernate.service systemd-suspend-then-hibernate.service

[Service]
Type=oneshot
TimeoutStartSec=300
ExecStart=/usr/bin/make_swapfile_to_hibernate

[Install]
RequiredBy=hibernate.target suspend-then-hibernate.target
